# image/supervisor_kafka.conf
[program:kafka]

command=/docker_start_kafka.sh

autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
